# .config/systemd/user/sway-autostart.target
[Unit]
Description=Sway session autostart
Documentation=man:systemd.special(7)
BindsTo=xdg-desktop-autostart.target graphical-session.target
