# .local/share/systemd/user/copyq.service
[Unit]
Description=CopyQ
BindsTo=graphical-session.target

[Service]
Type=simple
ExecStart=copyq

[Install]
WantedBy=xdg-desktop-autostart.target

